# rssd.service
[Unit]
Description=RSS polling daemon

[Service]
Type=oneshot
ExecStart=/usr/bin/rssd
